# Rules/99-power-profile.rules
SUBSYSTEM=="power_supply", ATTR{type}=="Mains", RUN+="/home/nico/Debian-i3-Setup/Scripts/SwitchPowerProfiles.sh"